# crontabs
#GTD
*/2 * * * * /home/philipp/gtd/move-new-todos.py INBOX INBOX.GTD.Input
*/2 * * * * /home/philipp/gtd/move-new-delayed-todos.py INBOX INBOX.GTD.Delayed



0 08 * * 1 /home/philipp/gtd/move-delayed-todos.py INBOX.GTD.Delayed INBOX.GTD.Today
0 08 * * 1 /home/philipp/gtd/moveall.py INBOX.GTD.Monday INBOX.GTD.Today
0 08 * * 2 /home/philipp/gtd/moveall.py INBOX.GTD.Tuesday INBOX.GTD.Today
0 08 * * 3 /home/philipp/gtd/moveall.py INBOX.GTD.Wednesday INBOX.GTD.Today
0 08 * * 4 /home/philipp/gtd/moveall.py INBOX.GTD.Thursday INBOX.GTD.Today
0 08 * * 5 /home/philipp/gtd/moveall.py INBOX.GTD.Friday INBOX.GTD.Today
0 08 * * 6 /home/philipp/gtd/moveall.py INBOX.GTD.Saturday INBOX.GTD.Today
0 08 * * 7 /home/philipp/gtd/moveall.py INBOX.GTD.Sunday INBOX.GTD.Today

1 08 * * * /home/philipp/gtd/moveall.py INBOX.GTD.Tomorrow INBOX.GTD.Today
4 08 * * * /home/philipp/gtd/unseenall.py INBOX.GTD.Today
2 08 * * 1 /home/philipp/gtd/moveall.py "INBOX.GTD.Next Week" INBOX.GTD.Today

3 08 * * * /home/philipp/gtd/copyall.py INBOX.GTD.Repeaters.Daily INBOX.GTD.Today
3 08 * * 1 /home/philipp/gtd/copyall.py INBOX.GTD.Repeaters.Weekly.Monday INBOX.GTD.Today
3 08 * * 2 /home/philipp/gtd/copyall.py INBOX.GTD.Repeaters.Weekly.Tuesday INBOX.GTD.Today
3 08 * * 3 /home/philipp/gtd/copyall.py INBOX.GTD.Repeaters.Weekly.Wednesday INBOX.GTD.Today
3 08 * * 4 /home/philipp/gtd/copyall.py INBOX.GTD.Repeaters.Weekly.Thursday INBOX.GTD.Today
3 08 * * 5 /home/philipp/gtd/copyall.py INBOX.GTD.Repeaters.Weekly.Friday INBOX.GTD.Today
3 08 * * 6 /home/philipp/gtd/copyall.py INBOX.GTD.Repeaters.Weekly.Saturday INBOX.GTD.Today
3 08 * * 7 /home/philipp/gtd/copyall.py INBOX.GTD.Repeaters.Weekly.Sunday INBOX.GTD.Today
